# deploy/systemd/rdda-xray.service
[Unit]
Description=RDDA xray-core data plane
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/xray run -c /etc/rdda/xray.json
Restart=on-failure
RestartSec=3
# Run as the dedicated rdda user (owns /etc/rdda, mode 0700) so xray can read
# its config. DynamicUser cannot — its transient UID is denied on the 0700 dir.
User=rdda
Group=rdda
AmbientCapabilities=CAP_NET_BIND_SERVICE
NoNewPrivileges=yes
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
