# deploy/systemd/rdda-update.timer
[Unit]
Description=RDDA self-update timer (opt-in; disabled by default)

[Timer]
OnBootSec=15min
OnUnitActiveSec=24h
# Stagger opted-in nodes so a bad release does not reach them all at once.
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target
