# deploy/systemd/rdda-nfqws.service
[Unit]
Description=RDDA nfqws2 egress DPI-desync (RU node, fail-open)
After=network-online.target rdda-singbox.service
Wants=network-online.target

[Service]
# Flags MUST match `rdda render nfqws` (profile + a fooling method). The
# --dpi-desync-fooling=badsum is load-bearing: without it the injected fake
# packet reaches the real server and corrupts every TLS handshake, including
# unblocked hosts and the RU→EU tunnel. Re-render if you change the profile.
ExecStartPre=/usr/sbin/nft -f /etc/rdda/rdda-nfqws.nft
ExecStart=/usr/local/bin/nfqws2 --qnum=200 --dpi-desync=fake,split2 --dpi-desync-fooling=badsum
ExecStopPost=-/usr/sbin/nft delete table inet rdda_nfqws
Restart=on-failure
RestartSec=3
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
