# deploy/systemd/rdda-alert.service
[Unit]
Description=RDDA EU alerting (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/rdda alert
User=rdda
Group=rdda
