# deploy/systemd/rdda-sub.service
[Unit]
Description=RDDA subscription server
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/rdda serve --addr 127.0.0.1:8080
Restart=on-failure
RestartSec=3
User=rdda
Group=rdda
ReadOnlyPaths=/etc/rdda
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
