# deploy/systemd/rdda-sub.service
[Unit]
Description=RDDA subscription server
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/rdda serve --addr 127.0.0.1:8080
Restart=on-failure
RestartSec=3
User=rdda
Group=rdda
# The sub server reads config/clients AND writes the last RU health beat
# (ru-health.json) here, so /etc/rdda must be writable to it. OS perms
# (0700, owned by rdda) still restrict the directory to this service.
ReadWritePaths=/etc/rdda
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
