# deploy/systemd/rdda-pull.service
[Unit]
Description=RDDA RU config pull-sync (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/rdda/pull.env
ExecStart=/usr/local/bin/rdda pull --from ${RDDA_PULL_FROM} --token ${RDDA_PULL_TOKEN}
User=rdda
Group=rdda
# No NoNewPrivileges: the default --reload-cmd is `sudo systemctl reload-or-restart
# rdda-singbox` (the rdda user is granted exactly that one command via
# /etc/sudoers.d/rdda-reload). NoNewPrivileges=yes would block sudo's setuid
# escalation, so the post-pull reload would fail and this oneshot would exit non-zero.
